# packages/snapcast/config/snapclient.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2020 OpenWrt.org

START=95
USE_PROCD=1

PROGARGS="-s pcm.default"

start_service() {
    if [ "$1" = "loopback" ]; then
        PROGARGS="$PROGARGS -h 127.0.0.1"
    fi

    procd_open_instance
    procd_set_param command /usr/bin/snapclient $PROGARGS
    procd_set_param respawn
    procd_append_param respawn 3600
    procd_append_param respawn 5
    procd_append_param respawn 5
    procd_close_instance
}
